# s6-sample/config/etc/services.d/inifinte-hello/run
#!/bin/sh

while true; do sleep 2; echo Hello; done